# sys_config/school_schedule.conf
[program:school_schedule]
command=/home/sean/school_schedule/sys_config/gunicorn_start
user=sean
stdout_logfile=/home/sean/school_schedule/logs/gunicorn_supervisor.log
redirect_stderr=true
