# sys_config/school_schedule.conf
[program:school_schedule]
command=~/school_schedule/sys_config/gunicorn_start
user=ubuntu
stdout_logfile=~/school_schedule/logs/gunicorn_supervisor.log
redirect_stderr=true
